# build/varnish/files/manifest-varnish.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
  Copyright 2011 Stephan Pietzko http://s2p.de/Sysadmin/SolarisVarnishManifest
  All rights reserved. Use is subject to license terms.
-->
<service_bundle type='manifest' name='varnish'>
  <service name='network/http' type='service' version='1'>
    <!--
      Because we may have multiple instances of network/http
      provided by different implementations, we keep dependencies
      and methods within the instance.
    -->
    <instance name='varnish' enabled='false'>
      <exec_method type='method' name='start' timeout_seconds='60' exec='/usr/local/bin/amd64/varnishd -P /var/run/varnishd.pid -f /etc/varnish/varnish.vcl -s malloc,3G -a 127.0.0.1:8080 -u varnish' />
      <exec_method type='method' name='stop' timeout_seconds='60' exec='/usr/bin/kill `cat /var/run/varnishd.pid`' />

      <!-- sub-process core dumps shouldn't restart session -->
      <property_group name='startd' type='framework'>
        <propval name='ignore_error' type='astring' value='core,signal' />
      </property_group>
    </instance>

    <template>
      <common_name>
        <loctext xml:lang='C'>Varnish 3 reverse proxy server</loctext>
      </common_name>
      <documentation>
        <manpage title='varnishd' section='8' manpath='/usr/local/share/man/man1/varnishd.1' />
        <doc_link name='varnish-cache.org' uri='https://www.varnish-cache.org/' />
      </documentation>
    </template>
  </service>
</service_bundle>