# build/consul/files/consul.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='consul'>
  <service name='monitoring/consul' type='service' version='1'>
    <create_default_instance enabled='true' />
    <single_instance />
    
    <dependency name='fs' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local' />
    </dependency>
    
    <dependency name='net' grouping='require_all' restart_on='none' type='service'>
       <service_fmri value='svc:/network/loopback' />
    </dependency>
    
    <dependency name='physical' grouping='optional_all' restart_on='none' type='service'>
      <service_fmri value='svc:/network/physical:default' />
    </dependency>
    
    <dependency name='config-file' grouping='require_all' restart_on='none' type='path'>
      <service_fmri value='file://localhost/etc/consul.d/config.json' />
    </dependency>
    
    <exec_method
       type='method'
       name='start'
       exec='/usr/local/bin/amd64/consul agent -config-dir /etc/consul.d >/var/log/consul.log 2>&1'
       timeout_seconds='-1'>
       <method_context>
         <method_credential user='consul' group='consul' />
       </method_context>
    </exec_method>
    
    <exec_method
      type='method'
      name='stop'
      exec=":kill"
      timeout_seconds='-1'>
    </exec_method>
    
    <exec_method
      type='method'
      name='reload'
      exec='/usr/local/bin/amd64/consul reload'
      timeout_seconds='-1'>
    </exec_method>
    
    <exec_method
      type='method'
      name='refresh'
      exec='/usr/local/bin/amd64/consul reload'
      timeout_seconds='-1'>
    </exec_method>
  
  </service>

</service_bundle>
